# vars/defaults/files
#!/bin/bash

remove_files=(
  '/boot/cmdline.txt'
  '/etc/init.d/resize2fs_once'
  '/etc/nv/nvfirstboot'
  '/etc/profile.d/armbian-check-first-login-reboot.sh'
  '/etc/profile.d/armbian-check-first-login.sh'
  '/etc/rc3.d/S01resize2fs_once'
  '/etc/ssh/ssh_host_*_key*'
  '/etc/systemd/system/default.target.wants/ssh.service'
  '/usr/bin/raspi-config'
  '/var/lib/systemd/deb-systemd-helper-enabled/multi-user.target.wants/NetworkManager.service'
)

null_link_files=(
  '/etc/init.d/dhcpcd'
  '/etc/systemd/network/99-default.link'
  '/etc/systemd/system/default.target.wants/avahi-daemon.service'
  '/etc/systemd/system/default.target.wants/dhcpcd.service'
  '/etc/systemd/system/default.target.wants/nfs-client.target'
  '/etc/systemd/system/default.target.wants/raspberrypi-net-mods.service'
  '/etc/systemd/system/default.target.wants/regenerate_ssh_host_keys.service'
  '/etc/systemd/system/default.target.wants/remote-fs.target'
  '/etc/systemd/system/default.target.wants/triggerhappy.service'
  '/etc/systemd/system/multi-user.target.wants/NetworkManager.service'
  '/etc/systemd/system/multi-user.target.wants/armbian-firstrun-config.service'
  '/lib/systemd/system/NetworkManager.service'
  '/lib/systemd/system/apply_noobs_os_config.service'
  '/lib/systemd/system/armbian-firstrun-config.service'
  '/lib/systemd/system/armbian-firstrun.service'
  '/lib/systemd/system/armbian-resize-filesystem.service'
  '/lib/systemd/system/dhcpcd.service'
  '/lib/systemd/system/raspberrypi-net-mods.service'
  '/lib/systemd/system/regenerate_ssh_host_keys.service'
  '/usr/lib/armbian/armbian-firstrun'
  '/usr/lib/armbian/armbian-firstrun-config'
  '/usr/local/bin/raspi-config'
  '/var/lib/systemd/deb-systemd-helper-enabled/multi-user.target.wants/regenerate_ssh_host_keys.service'
)
